# systemd/hyprlaunch.service
[Service]
Type=oneshot
ExecStart=%h/.local/bin/hyprlaunch

[Install]
WantedBy=default.target
